# deb/borhan-nginx/debian/prerm
#!/bin/sh
if [ "$1" = "remove" -o $1 = "purge" ] ; then
    /usr/sbin/service borhan-nginx stop >/dev/null 2>&1
    update-rc.d -f borhan-nginx remove
fi
